# supervisord.conf
[supervisord]
nodaemon=true

[program:dnsmasq]
command = /usr/sbin/dnsmasq -C /etc/dnsmasq.conf --no-daemon
autostart = true
autorestart = true
